# build/farming_simulator.conf
[program:fs25_server]
command=/usr/local/bin/autostart_fs25.sh
user=nobody
environment=HOME="/home/nobody",USER="nobody",DISPLAY=":0",XAUTHORITY="/home/nobody/.Xauthority"
